# scripts/eink-panel.service
[Unit]
Description=E-Ink Control Panel
After=network.target

[Service]
Type=simple
User=tjoskar
WorkingDirectory=/home/tjoskar/control-panel
ExecStart=/home/tjoskar/.local/share/fnm/aliases/default/bin/node --env-file=.env --enable-source-maps dist/home-dashboard/main.js
Restart=on-failure
RestartSec=10

# Environment
Environment=NODE_ENV=production

# Logging
StandardOutput=append:/home/tjoskar/control-panel/logs/eink-panel.log
StandardError=append:/home/tjoskar/control-panel/logs/eink-panel.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/tjoskar/control-panel/logs /tmp/eink-panel

[Install]
WantedBy=multi-user.target
